# hyperglass-ansible/roles/hyperglass/templates/hyperglass.service.j2
[Unit]
Description=hyperglass
After=network.target

[Service]
ExecStart=/usr/local/bin/hyperglass start
Restart=always
User=root

[Install]
WantedBy=multi-user.target
